# deploy/systemd/mes-pilot.service
[Unit]
Description=MES Pilot Preview
After=network-online.target mes-pilot-credential-rotation-recovery.service
Wants=network-online.target
Requires=mes-pilot-credential-rotation-recovery.service

[Service]
Type=simple
User=mes-pilot
Group=mes-pilot
SupplementaryGroups=mes-pilot-data
WorkingDirectory=/srv/mes/pilot/app
EnvironmentFile=/etc/mes/mes-pilot.env
Environment=HOME=/srv/mes/pilot/runtime
Environment=NPM_CONFIG_CACHE=/srv/mes/pilot/runtime/npm-cache
ExecStartPre=/usr/bin/npm run server:preflight
ExecStart=/usr/bin/npm run preview
Restart=always
RestartSec=5
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
KeyringMode=private
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
LimitCORE=0
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/srv/mes/pilot/app /srv/mes/pilot/releases
ReadWritePaths=/srv/mes/pilot/shared-state /srv/mes/pilot/backups /srv/mes/pilot/audit /srv/mes/pilot/runtime

[Install]
WantedBy=multi-user.target
