# deploy/systemd/mes-user-testing.service
[Unit]
Description=MES prototype user-testing contour
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mes
Group=mes
WorkingDirectory=/opt/mes/current
EnvironmentFile=/etc/mes/mes-user-testing.env
ExecStartPre=/usr/bin/npm run server:preflight
ExecStart=/usr/bin/npm run preview
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/mes/user-testing /opt/mes/current

[Install]
WantedBy=multi-user.target
