# deploy/systemd/mes-user-testing.service
[Unit]
Description=MES User Testing Preview
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/srv/mes/user-testing/app
EnvironmentFile=/etc/mes/mes-user-testing.env
ExecStartPre=/usr/bin/npm run server:preflight
ExecStart=/usr/bin/npm run preview
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/mes/user-testing

[Install]
WantedBy=multi-user.target
