# ops/security/mes-pilot-domain-runtime-credential-check.service
[Unit]
Description=MES Pilot runtime PostgreSQL credential check
After=postgresql.service

[Service]
Type=oneshot
User=mes-pilot
Group=mes-pilot
WorkingDirectory=/srv/mes/pilot/app
EnvironmentFile=/etc/mes/mes-pilot-domain.env
Environment=HOME=/nonexistent
ExecStart=/usr/local/libexec/mes/pilot-runtime-security-dispatch check-postgres-credential.mjs --variable=DATABASE_URL --expected-role=mes_app
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
LimitCORE=0
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/srv/mes/pilot/app /srv/mes/pilot/releases /usr/local/libexec/mes
