# ops/security/mes-pilot-credential-rotation-recovery.service
[Unit]
Description=Recover interrupted MES Pilot credential rotation
After=local-fs.target postgresql.service
Before=mes-pilot.service mes-pilot-domain-migrate.service mes-pilot-domain-import.service mes-pilot-domain-snapshot-sync.service

[Service]
Type=oneshot
User=root
Group=root
# Keep both recovery passes inside one canonical fd9 ownership interval. During
# an intentional app verification the fixed wrapper admits only the exact live
# release/runtime intent and exits without executing either recovery pass.
ExecStart=/bin/bash /usr/local/libexec/mes/active-bundle/with-pilot-release-authority-lock.sh --operation=runtime-security-recovery --busy-policy=app-intent -- /bin/bash -ceu '/usr/local/libexec/mes/pilot-runtime-security-dispatch recover-pilot-credential-rotation.sh; /usr/local/libexec/mes/pilot-runtime-security-dispatch recover-pilot-uid-cutover.sh'
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_FOWNER CAP_SETGID CAP_SETUID
AmbientCapabilities=
LimitCORE=0
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
