# ops/postgres/mes-pilot-domain-snapshot-sync.timer
[Unit]
Description=Retry MES Pilot PostgreSQL compatibility outbox

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
RandomizedDelaySec=3s
Persistent=true
Unit=mes-pilot-domain-snapshot-sync.service

[Install]
WantedBy=timers.target
