# ops/postgres/mes-pilot-domain-migrate.service
[Unit]
Description=MES Pilot domain database migration
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/srv/mes/pilot/app
EnvironmentFile=/etc/mes/mes-pilot-domain.env
ExecStart=/usr/bin/npm run domain:postgres:migrate
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/mes/pilot

[Install]
WantedBy=multi-user.target
