# ops/postgres/mes-pilot-domain-migrate.service
[Unit]
Description=MES Pilot domain database migration
After=network-online.target postgresql.service
Wants=network-online.target
Requires=mes-pilot-credential-rotation-recovery.service
After=mes-pilot-credential-rotation-recovery.service

[Service]
Type=oneshot
User=mes-pilot-migrator
Group=mes-pilot-migrator
SupplementaryGroups=mes-pilot-data
WorkingDirectory=/srv/mes/pilot/app
EnvironmentFile=/etc/mes/mes-pilot-domain-migrator.env
Environment=HOME=/nonexistent
ExecStart=/usr/bin/node /srv/mes/pilot/app/scripts/domain-postgres-migrate.mjs --schema-only
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
KeyringMode=private
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
LimitCORE=0
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/srv/mes/pilot/app /srv/mes/pilot/releases

[Install]
WantedBy=multi-user.target
